# supervisord.conf
[supervisord]
nodaemon=true

[program:kestrel]
directory=/app/TreintijdenDash
command=k kestrel